# .config/systemd/user/metapac-unsorted.service
[Unit]
Description=Snapshot unmanaged Metapac packages to unsorted.toml

[Service]
Type=oneshot
ExecStart=%h/.local/bin/metapac-periodic-sync.sh

[Install]
WantedBy=default.target
